# crontab

* * * * * cd /home/airinsider/var/www/html/airinsider && php artisan schedule:run >> /dev/null 2>&1
