# crontab

* * * * * cd /home/test/laravel-dev/airinsider && php artisan schedule:run >> /dev/null 2>&1
